# scripts/systemd/ifdb-uwsgi-zok.service
[Unit]
Description=Ifdb Zok Uwsgi server
After=syslog.target

[Service]
UMask=0002
ExecStart=/home/ifdb/venv/bin/uwsgi --ini /home/ifdb/configs/uwsgi-zok.ini
Restart=always
KillSignal=SIGQUIT

[Install]
WantedBy=multi-user.target